# kyberguard-ioc-watcher.timer
[Unit]
Description=KyberGuard IOC-Watcher Timer — täglich 02:30 UTC

[Timer]
OnCalendar=*-*-* 02:30:00 UTC
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
